# ghostmode/temp_dev_files/old_versions/ghost_idlewatch.service
[Unit]
Description=GhostMode Self-Wipe Timer
After=network.target

[Service]
ExecStart=%h/.local/bin/ghost_idlewatch.sh
Restart=on-failure

[Install]
WantedBy=default.target
